# dockerfiles/novnc-ubuntu-base/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/markus

[program:tigervnc]
command=vncserver :1 -SecurityTypes None
user=markus
directory=/home/markus